# debian/prerm
#!/bin/bash
# PowerVigil™ Pre-Removal Script

set -e

case "$1" in
    remove|upgrade|deconfigure)
        echo "PowerVigil™: Stopping services..."
        
        # Stop services if they're running
        systemctl stop powervigil-watchdog.service 2>/dev/null || true
        systemctl stop powervigil.service 2>/dev/null || true
        
        # Disable services
        systemctl disable powervigil-watchdog.service 2>/dev/null || true
        systemctl disable powervigil.service 2>/dev/null || true
        ;;
        
    failed-upgrade)
        ;;
        
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
